# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid ivialivia
chdir /home/ivialivia/sites/SITENAME/sourcei

exec ../virtualenv/bin/gunicorn \
	--bind unix:/tmp/SITENAME.scoket \
	superlists.wsgi:application
